# mqttserver/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mqttserver]
command=/usr/local/sbin/mosquitto